# cloud_setup/webxes-health-monitor.service
[Unit]
Description=WEBXES Tech Health Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ai_employee_vault
ExecStart=/usr/bin/python3 cloud_setup/health_monitor.py
Restart=on-failure
RestartSec=60
EnvironmentFile=/opt/ai_employee_vault/.env
StandardOutput=append:/opt/ai_employee_vault/Logs/health_monitor_systemd.log
StandardError=append:/opt/ai_employee_vault/Logs/health_monitor_systemd.log

[Install]
WantedBy=multi-user.target
