# cloud_setup/webxes-orchestrator.service
[Unit]
Description=WEBXES Tech Orchestrator (Scheduler)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/ai_employee_vault
ExecStart=/usr/bin/python3 orchestrator.py
Restart=on-failure
RestartSec=30
EnvironmentFile=/opt/ai_employee_vault/.env
StandardOutput=append:/opt/ai_employee_vault/Logs/orchestrator_systemd.log
StandardError=append:/opt/ai_employee_vault/Logs/orchestrator_systemd.log

[Install]
WantedBy=multi-user.target
